# scripts/com.ttcguide.daily-blog-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ttcguide.daily-blog-publish</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/satoutakuma/Desktop/webpage.new/scripts/daily-blog-publish.sh</string>
    </array>

    <!-- 毎日 09:00 に実行（research:deep の1時間後） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/satoutakuma/Desktop/webpage.new/scripts/logs/launchd-blog-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/satoutakuma/Desktop/webpage.new/scripts/logs/launchd-blog-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/satoutakuma/Desktop/webpage.new</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
